# systemd/user/ensure-dirs.service
[Unit]

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/bin/ensure-dirs

[Install]
WantedBy=multi-user.target
