# systemd/user/xkb-load.service
[Unit]

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/bin/xkb-load

[Install]
WantedBy=multi-user.target
