# systemd/user/xkb-monitor.service
[Unit]

[Service]
Type=forking
ExecStart=%h/.local/bin/xkb-monitor

[Install]
WantedBy=multi-user.target
